# systemd/sb-live.service
[Unit]
Description=SB Watchbot ICT Silver Bullet live
After=network-online.target sb-live-fetch.service
Wants=sb-live-fetch.service

[Service]
Type=simple
WorkingDirectory=/opt/sb-watchbot
EnvironmentFile=/etc/sb-watchbot/env
Environment=HOME=/opt/sb-watchbot
# Use module CLI (supports --poll). Tee logs for visibility.
ExecStart=/bin/bash -c 'source .venv/bin/activate && PYTHONPATH=. python -m sbwatch.app.live_sb --csv /opt/sb-watchbot/live/nq_1m.csv --poll 5 --heartbeat --daily-pings 2>&1 | tee -a out/live.log'
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
