# systemd/sb-live.service
[Unit]
Description=SB Simple - Live Alerts (Databento Live)
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sb-watchbot.env
Environment=PYTHONPATH=/opt/sb-simple/src
WorkingDirectory=/opt/sb-simple
ExecStart=/opt/sb-simple/.venv/bin/python /opt/sb-simple/run_live.py
RuntimeMaxSec=3720
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
