# systemd/sb-live-fetch.service
[Unit]
Description=SB Watchbot live data fetch (Databento 1m)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sb-watchbot
EnvironmentFile=/etc/sb-watchbot/env
Environment=HOME=/opt/sb-watchbot
ExecStart=/bin/bash -c 'bash scripts/live_fetch_nq.sh'
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
